# crontab
# update all containers
0 0 * * * (cd ~/; make update) > ~/crontab_log
